# supervisord.conf
[supervisord]
nodaemon=true

[program:scout_realtime]
command=/bin/bash -c "scout_realtime"

[program:ragios]
command=/bin/bash -c "./ragios s start"
